# packages/bird-openwrt/bird4-openwrt/src/init.d/bird4
#!/bin/sh /etc/rc.common

# Copyright (C) 2014-2017 - Eloi Carbo
#
# This program is free software: you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation, either version 3 of the License, or
# (at your option) any later version.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program.  If not, see <http://www.gnu.org/licenses/>.
#

BIRD="bird4"
BIRD_CONFIG="/etc/$BIRD.conf"

START=99
STOP=10

SERVICE_DAEMONIZE=1
SERVICE_USE_PID=1
SERVICE_PID_FILE="/var/run/$BIRD.pid"

BIRD_BIN="/usr/sbin/$BIRD"

. /etc/${BIRD}/init.d/${BIRD}-lib.sh

start() {
    config_load bird4
    local use_UCI_config
    get use_UCI_config 'bird'
	
    if [ -z "$use_UCI_config" -o "$use_UCI_config" = "0" ]; then
        service_start $BIRD_BIN -d -c $BIRD_CONFIG -P $SERVICE_PID_FILE
    else
        #Set Bird4 configuration location:
        local UCI_config_File
        get UCI_config_File 'bird'
        BIRD_CONFIG="${UCI_config_file}:-/tmp/bird4.conf"
        #Setup the basic configuration
        prepare_global 'global'
        config_foreach prepare_kernel 'kernel'
        config_foreach prepare_static 'static'
        config_foreach prepare_device 'device'
        #config_foreach prepare_direct 'direct'
		
        #Setup protocol`s configuration: BGP
        config_foreach prepare_bgp_template 'bgp_template'
        config_foreach prepare_bgp 'bgp'
        config_foreach prepare_bgp_filters 'filter'

        #Setup protocol`s configuration: OSPF
	    config_foreach prepare_ospf_instance 'ospf'
        config_foreach prepare_ospf_filters 'filter'

        #Start the service
        service_start $BIRD_BIN -d -c $BIRD_CONFIG -P $SERVICE_PID_FILE
    fi
}

stop() {
    service_stop $BIRD_BIN
}

restart() {
    stop
    start
}


reload() {
    service_reload $BIRD_BIN
}
